# deploy/systemd/hedgekit-dashboard.service
[Unit]
Description=hedgekit dashboard process (SPEC Process D)
After=network.target

[Service]
ExecStart=/usr/bin/env hedgekit run --process dashboard
Restart=on-failure

[Install]
WantedBy=multi-user.target
