# deploy/systemd/hedgekit-pipeline.service
[Unit]
Description=hedgekit pipeline process (SPEC Process A)
After=network.target

[Service]
ExecStart=/usr/bin/env hedgekit run --process pipeline
Restart=on-failure

[Install]
WantedBy=multi-user.target
