# deploy/systemd/hedgekit-order-gateway.service
[Unit]
Description=hedgekit order gateway process (SPEC Process C)
After=network.target

[Service]
ExecStart=/usr/bin/env hedgekit run --process order_gateway
Restart=on-failure

[Install]
WantedBy=multi-user.target
