# deploy/systemd/hedgekit-riskkernel.service
[Unit]
Description=hedgekit risk kernel process (SPEC Process B)
After=network.target

[Service]
ExecStart=/usr/bin/env hedgekit run --process riskkernel
Restart=on-failure

[Install]
WantedBy=multi-user.target
